# packages/auth.regarde.dev/deploy/config/systemd.service
# Systemd service template for auth.regarde.dev
# This file will be processed by deployment scripts to replace template variables
# with actual values from the .env configuration

[Unit]
Description=auth.regarde.dev Authentication Service
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_ROOT}}
ExecStart=/usr/bin/pnpm --filter @regarde-dev/auth.regarde.dev start:bun
Environment=PATH={{BUN_PATH}}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV={{NODE_ENV}}
Environment=PORT=3000
Environment=APP_PUBLIC_HOSTNAME=auth.regarde.dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
